# root/etc/systemd/system/ipsec.service.d/nethserver.conf
[Service]
# on ipsec stop delete all conntrack connections
ExecStopPost=/usr/libexec/nethserver/nethserver-ipsec-tunnels-delete-conntrack
